# bin/testing/update_carmusr/update_carmusr.cron
# anfa
#SHELL = /bin/tcsh
#25 20 * * * /carm/proj/skcms/CARMUSR/sk_cms_user/bin/desmond.sh stop |& tee -a /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr.log2
#30 03 * * * /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr |& tee -a /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr.log2 | mail -s 'Update of sk_cms_main_trunk user' gronlund
25 19 * * * /carm/proj/skcms/CARMUSR/sk_cms_user/bin/desmond.sh stop 2>&1 | tee -a /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr.log2 >/dev/null
30 03 * * * /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr 2>&1 | tee /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr.log | tee -a /carm/proj/skcms/CARMUSR/sk_cms_user/update_carmusr.log2 | mail -s 'Update of sk_cms_main_trunk user' gronlund
